# aiobi-update/systemd/user/aiobi-update-notify.path
[Unit]
Description=Aïobi OS — watch for update-available notify trigger

[Path]
PathModified=/run/aiobi-update/notify

[Install]
WantedBy=default.target
